# package/pifii/udpreport/rootfs/etc/rc.d/S93homemac
#!/bin/sh /etc/rc.common
START=93
EXTRA_COMMANDS=""

start () {
        /usr/sbin/homemac&
}

stop () {
        service_stop /usr/sbin/homemac
}


